# services/spark-shell-service/run
#!/bin/bash
exec spark-submit --conf spark.driver.userClassPathFirst=true --conf spark.sql.shuffle.partitions=200 --class com.thinkbiganalytics.spark.SparkShellServer --jars "/root/data-lake-accelerator/core/db-schema/db-schema-model/target/thinkbig-db-schema-model-0.1.0-SNAPSHOT.jar,/root/data-lake-accelerator/services/feed-manager-service/feed-manager-api/target/thinkbig-feed-manager-service-0.1.0-SNAPSHOT.jar,/root/.m2/repository/com/google/guava/guava/18.0/guava-18.0.jar,/root/.m2/repository/org/glassfish/jersey/containers/jersey-container-jetty-http/2.22.2/jersey-container-jetty-http-2.22.2.jar,/root/.m2/repository/org/glassfish/hk2/external/javax.inject/2.4.0-b34/javax.inject-2.4.0-b34.jar,/root/.m2/repository/org/eclipse/jetty/jetty-server/9.1.1.v20140108/jetty-server-9.1.1.v20140108.jar,/root/.m2/repository/javax/servlet/javax.servlet-api/3.1.0/javax.servlet-api-3.1.0.jar,/root/.m2/repository/org/eclipse/jetty/jetty-http/9.1.1.v20140108/jetty-http-9.1.1.v20140108.jar,/root/.m2/repository/org/eclipse/jetty/jetty-util/9.1.1.v20140108/jetty-util-9.1.1.v20140108.jar,/root/.m2/repository/org/eclipse/jetty/jetty-io/9.1.1.v20140108/jetty-io-9.1.1.v20140108.jar,/root/.m2/repository/org/eclipse/jetty/jetty-continuation/9.1.1.v20140108/jetty-continuation-9.1.1.v20140108.jar,/root/.m2/repository/org/glassfish/jersey/core/jersey-common/2.22.2/jersey-common-2.22.2.jar,/root/.m2/repository/javax/ws/rs/javax.ws.rs-api/2.0.1/javax.ws.rs-api-2.0.1.jar,/root/.m2/repository/javax/annotation/javax.annotation-api/1.2/javax.annotation-api-1.2.jar,/root/.m2/repository/org/glassfish/jersey/bundles/repackaged/jersey-guava/2.22.2/jersey-guava-2.22.2.jar,/root/.m2/repository/org/glassfish/hk2/hk2-api/2.4.0-b34/hk2-api-2.4.0-b34.jar,/root/.m2/repository/org/glassfish/hk2/hk2-utils/2.4.0-b34/hk2-utils-2.4.0-b34.jar,/root/.m2/repository/org/glassfish/hk2/external/aopalliance-repackaged/2.4.0-b34/aopalliance-repackaged-2.4.0-b34.jar,/root/.m2/repository/org/glassfish/hk2/hk2-locator/2.4.0-b34/hk2-locator-2.4.0-b34.jar,/root/.m2/repository/org/javassist/javassist/3.18.1-GA/javassist-3.18.1-GA.jar,/root/.m2/repository/org/glassfish/hk2/osgi-resource-locator/1.0.1/osgi-resource-locator-1.0.1.jar,/root/.m2/repository/org/glassfish/jersey/core/jersey-server/2.22.2/jersey-server-2.22.2.jar,/root/.m2/repository/org/glassfish/jersey/core/jersey-client/2.22.2/jersey-client-2.22.2.jar,/root/.m2/repository/org/glassfish/jersey/media/jersey-media-jaxb/2.22.2/jersey-media-jaxb-2.22.2.jar,/root/.m2/repository/javax/validation/validation-api/1.1.0.Final/validation-api-1.1.0.Final.jar,/root/.m2/repository/org/glassfish/jersey/media/jersey-media-json-jackson/2.22.2/jersey-media-json-jackson-2.22.2.jar,/root/.m2/repository/org/glassfish/jersey/ext/jersey-entity-filtering/2.22.2/jersey-entity-filtering-2.22.2.jar,/root/.m2/repository/com/fasterxml/jackson/core/jackson-annotations/2.6.2/jackson-annotations-2.6.2.jar,/root/.m2/repository/commons-logging/commons-logging/1.1.3/commons-logging-1.1.3.jar,/root/.m2/repository/log4j/log4j/1.2.17/log4j-1.2.17.jar,/root/.m2/repository/org/slf4j/slf4j-api/1.7.12/slf4j-api-1.7.12.jar,/root/.m2/repository/org/slf4j/slf4j-log4j12/1.7.10/slf4j-log4j12-1.7.10.jar,/root/.m2/repository/javax/inject/javax.inject/1/javax.inject-1.jar,/root/.m2/repository/aopalliance/aopalliance/1.0/aopalliance-1.0.jar,/root/.m2/repository/org/apache/commons/commons-lang3/3.4/commons-lang3-3.4.jar,/root/.m2/repository/com/fasterxml/jackson/core/jackson-databind/2.6.2/jackson-databind-2.6.2.jar,/root/.m2/repository/com/fasterxml/jackson/core/jackson-core/2.6.2/jackson-core-2.6.2.jar,/root/.m2/repository/io/swagger/swagger-jersey2-jaxrs/1.5.8/swagger-jersey2-jaxrs-1.5.8.jar,/root/.m2/repository/io/swagger/swagger-jaxrs/1.5.8/swagger-jaxrs-1.5.8.jar,/root/.m2/repository/com/fasterxml/jackson/dataformat/jackson-dataformat-yaml/2.4.5/jackson-dataformat-yaml-2.4.5.jar,/root/.m2/repository/org/yaml/snakeyaml/1.12/snakeyaml-1.12.jar,/root/.m2/repository/com/fasterxml/jackson/dataformat/jackson-dataformat-xml/2.6.2/jackson-dataformat-xml-2.6.2.jar,/root/.m2/repository/com/fasterxml/jackson/module/jackson-module-jaxb-annotations/2.6.2/jackson-module-jaxb-annotations-2.6.2.jar,/root/.m2/repository/org/codehaus/woodstox/stax2-api/3.1.4/stax2-api-3.1.4.jar,/root/.m2/repository/io/swagger/swagger-core/1.5.8/swagger-core-1.5.8.jar,/root/.m2/repository/com/fasterxml/jackson/datatype/jackson-datatype-joda/2.6.2/jackson-datatype-joda-2.6.2.jar,/root/.m2/repository/joda-time/joda-time/2.9.2/joda-time-2.9.2.jar,/root/.m2/repository/io/swagger/swagger-models/1.5.8/swagger-models-1.5.8.jar,/root/.m2/repository/io/swagger/swagger-annotations/1.5.8/swagger-annotations-1.5.8.jar,/root/.m2/repository/org/reflections/reflections/0.9.10/reflections-0.9.10.jar,/root/.m2/repository/com/google/code/findbugs/annotations/2.0.1/annotations-2.0.1.jar,/root/.m2/repository/com/fasterxml/jackson/jaxrs/jackson-jaxrs-json-provider/2.6.2/jackson-jaxrs-json-provider-2.6.2.jar,/root/.m2/repository/com/fasterxml/jackson/jaxrs/jackson-jaxrs-base/2.6.2/jackson-jaxrs-base-2.6.2.jar,/root/.m2/repository/org/glassfish/jersey/containers/jersey-container-servlet-core/2.22.2/jersey-container-servlet-core-2.22.2.jar,/root/.m2/repository/org/glassfish/jersey/media/jersey-media-multipart/2.22.2/jersey-media-multipart-2.22.2.jar,/root/.m2/repository/org/jvnet/mimepull/mimepull/1.9.6/mimepull-1.9.6.jar,/root/.m2/repository/org/springframework/spring-context-support/4.2.5.RELEASE/spring-context-support-4.2.5.RELEASE.jar,/root/.m2/repository/org/springframework/spring-core/4.2.5.RELEASE/spring-core-4.2.5.RELEASE.jar,/root/.m2/repository/org/springframework/spring-beans/4.2.5.RELEASE/spring-beans-4.2.5.RELEASE.jar,/root/.m2/repository/org/springframework/spring-context/4.2.5.RELEASE/spring-context-4.2.5.RELEASE.jar,/root/.m2/repository/org/springframework/spring-aop/4.2.5.RELEASE/spring-aop-4.2.5.RELEASE.jar" --driver-java-options "-Xdebug -agentlib:jdwp=transport=dt_socket,server=y,suspend=n,address=5005 -Dcom.sun.management.jmxremote -Dcom.sun.management.jmxremote.port=9999 -Dcom.sun.management.jmxremote.authenticate=false -Dcom.sun.management.jmxremote.ssl=false -Dcom.sun.management.jmxremote.local.only=false -Djava.rmi.server.hostname=127.0.0.1" /root/data-lake-accelerator/services/spark-shell-service/target/thinkbig-spark-shell-service-0.1.0-SNAPSHOT.jar
