# services/spark-shell-service/run
#!/bin/bash
exec spark-submit --master local[4] --conf spark.scheduler.mode=FAIR --conf spark.driver.userClassPathFirst=true --conf spark.sql.shuffle.partitions=200 --class com.thinkbiganalytics.spark.SparkShellApp --driver-class-path /opt/thinkbig/thinkbig-spark-shell/conf --jars "/root/.m2/repository/org/glassfish/jersey/containers/jersey-container-servlet/2.22.2/jersey-container-servlet-2.22.2.jar,/root/.m2/repository/org/glassfish/jersey/containers/jersey-container-servlet-core/2.22.2/jersey-container-servlet-core-2.22.2.jar,/root/.m2/repository/org/glassfish/hk2/external/javax.inject/2.4.0-b34/javax.inject-2.4.0-b34.jar,/root/.m2/repository/org/glassfish/jersey/core/jersey-common/2.22.2/jersey-common-2.22.2.jar,/root/.m2/repository/javax/annotation/javax.annotation-api/1.2/javax.annotation-api-1.2.jar,/root/.m2/repository/org/glassfish/jersey/bundles/repackaged/jersey-guava/2.22.2/jersey-guava-2.22.2.jar,/root/.m2/repository/org/glassfish/hk2/hk2-api/2.4.0-b34/hk2-api-2.4.0-b34.jar,/root/.m2/repository/org/glassfish/hk2/hk2-utils/2.4.0-b34/hk2-utils-2.4.0-b34.jar,/root/.m2/repository/org/glassfish/hk2/external/aopalliance-repackaged/2.4.0-b34/aopalliance-repackaged-2.4.0-b34.jar,/root/.m2/repository/org/glassfish/hk2/hk2-locator/2.4.0-b34/hk2-locator-2.4.0-b34.jar,/root/.m2/repository/org/javassist/javassist/3.18.1-GA/javassist-3.18.1-GA.jar,/root/.m2/repository/org/glassfish/hk2/osgi-resource-locator/1.0.1/osgi-resource-locator-1.0.1.jar,/root/.m2/repository/org/glassfish/jersey/core/jersey-server/2.22.2/jersey-server-2.22.2.jar,/root/.m2/repository/org/glassfish/jersey/core/jersey-client/2.22.2/jersey-client-2.22.2.jar,/root/.m2/repository/org/glassfish/jersey/media/jersey-media-jaxb/2.22.2/jersey-media-jaxb-2.22.2.jar,/root/.m2/repository/javax/validation/validation-api/1.1.0.Final/validation-api-1.1.0.Final.jar,/root/.m2/repository/javax/ws/rs/javax.ws.rs-api/2.0.1/javax.ws.rs-api-2.0.1.jar,/root/.m2/repository/org/glassfish/jersey/media/jersey-media-json-jackson/2.22.2/jersey-media-json-jackson-2.22.2.jar,/root/.m2/repository/org/glassfish/jersey/ext/jersey-entity-filtering/2.22.2/jersey-entity-filtering-2.22.2.jar,/root/.m2/repository/com/fasterxml/jackson/core/jackson-annotations/2.6.2/jackson-annotations-2.6.2.jar,/root/.m2/repository/com/fasterxml/jackson/core/jackson-databind/2.6.2/jackson-databind-2.6.2.jar,/root/.m2/repository/com/fasterxml/jackson/core/jackson-core/2.6.2/jackson-core-2.6.2.jar,/root/.m2/repository/commons-logging/commons-logging/1.1.3/commons-logging-1.1.3.jar,/root/.m2/repository/joda-time/joda-time/2.9.2/joda-time-2.9.2.jar,/root/.m2/repository/org/springframework/boot/spring-boot-starter-jersey/1.3.3.RELEASE/spring-boot-starter-jersey-1.3.3.RELEASE.jar,/root/.m2/repository/org/springframework/boot/spring-boot-starter/1.3.3.RELEASE/spring-boot-starter-1.3.3.RELEASE.jar,/root/.m2/repository/org/springframework/boot/spring-boot/1.3.3.RELEASE/spring-boot-1.3.3.RELEASE.jar,/root/.m2/repository/org/springframework/boot/spring-boot-autoconfigure/1.3.3.RELEASE/spring-boot-autoconfigure-1.3.3.RELEASE.jar,/root/.m2/repository/org/yaml/snakeyaml/1.16/snakeyaml-1.16.jar,/root/.m2/repository/org/springframework/boot/spring-boot-starter-tomcat/1.3.3.RELEASE/spring-boot-starter-tomcat-1.3.3.RELEASE.jar,/root/.m2/repository/org/apache/tomcat/embed/tomcat-embed-core/8.0.32/tomcat-embed-core-8.0.32.jar,/root/.m2/repository/org/apache/tomcat/embed/tomcat-embed-el/8.0.32/tomcat-embed-el-8.0.32.jar,/root/.m2/repository/org/apache/tomcat/embed/tomcat-embed-logging-juli/8.0.32/tomcat-embed-logging-juli-8.0.32.jar,/root/.m2/repository/org/apache/tomcat/embed/tomcat-embed-websocket/8.0.32/tomcat-embed-websocket-8.0.32.jar,/root/.m2/repository/org/springframework/boot/spring-boot-starter-validation/1.3.3.RELEASE/spring-boot-starter-validation-1.3.3.RELEASE.jar,/root/.m2/repository/org/hibernate/hibernate-validator/5.2.4.Final/hibernate-validator-5.2.4.Final.jar,/root/.m2/repository/org/jboss/logging/jboss-logging/3.2.1.Final/jboss-logging-3.2.1.Final.jar,/root/.m2/repository/com/fasterxml/classmate/1.1.0/classmate-1.1.0.jar,/root/.m2/repository/org/springframework/spring-web/4.2.5.RELEASE/spring-web-4.2.5.RELEASE.jar,/root/.m2/repository/org/springframework/spring-aop/4.2.5.RELEASE/spring-aop-4.2.5.RELEASE.jar,/root/.m2/repository/aopalliance/aopalliance/1.0/aopalliance-1.0.jar,/root/.m2/repository/org/springframework/spring-beans/4.2.5.RELEASE/spring-beans-4.2.5.RELEASE.jar,/root/.m2/repository/org/springframework/spring-context/4.2.5.RELEASE/spring-context-4.2.5.RELEASE.jar,/root/.m2/repository/org/springframework/spring-expression/4.2.5.RELEASE/spring-expression-4.2.5.RELEASE.jar,/root/.m2/repository/org/glassfish/jersey/ext/jersey-bean-validation/2.22.2/jersey-bean-validation-2.22.2.jar,/root/.m2/repository/org/glassfish/jersey/ext/jersey-spring3/2.22.2/jersey-spring3-2.22.2.jar,/root/.m2/repository/org/glassfish/hk2/hk2/2.4.0-b34/hk2-2.4.0-b34.jar,/root/.m2/repository/org/glassfish/hk2/config-types/2.4.0-b34/config-types-2.4.0-b34.jar,/root/.m2/repository/org/glassfish/hk2/hk2-core/2.4.0-b34/hk2-core-2.4.0-b34.jar,/root/.m2/repository/org/glassfish/hk2/hk2-config/2.4.0-b34/hk2-config-2.4.0-b34.jar,/root/.m2/repository/org/jvnet/tiger-types/1.4/tiger-types-1.4.jar,/root/.m2/repository/org/glassfish/hk2/hk2-runlevel/2.4.0-b34/hk2-runlevel-2.4.0-b34.jar,/root/.m2/repository/org/glassfish/hk2/class-model/2.4.0-b34/class-model-2.4.0-b34.jar,/root/.m2/repository/org/glassfish/hk2/external/asm-all-repackaged/2.4.0-b34/asm-all-repackaged-2.4.0-b34.jar,/root/.m2/repository/org/glassfish/hk2/spring-bridge/2.4.0-b34/spring-bridge-2.4.0-b34.jar,/root/.m2/repository/io/swagger/swagger-jersey2-jaxrs/1.5.8/swagger-jersey2-jaxrs-1.5.8.jar,/root/.m2/repository/io/swagger/swagger-jaxrs/1.5.8/swagger-jaxrs-1.5.8.jar,/root/.m2/repository/com/fasterxml/jackson/dataformat/jackson-dataformat-yaml/2.4.5/jackson-dataformat-yaml-2.4.5.jar,/root/.m2/repository/com/fasterxml/jackson/dataformat/jackson-dataformat-xml/2.6.2/jackson-dataformat-xml-2.6.2.jar,/root/.m2/repository/com/fasterxml/jackson/module/jackson-module-jaxb-annotations/2.6.2/jackson-module-jaxb-annotations-2.6.2.jar,/root/.m2/repository/org/codehaus/woodstox/stax2-api/3.1.4/stax2-api-3.1.4.jar,/root/.m2/repository/io/swagger/swagger-core/1.5.8/swagger-core-1.5.8.jar,/root/.m2/repository/com/fasterxml/jackson/datatype/jackson-datatype-joda/2.6.2/jackson-datatype-joda-2.6.2.jar,/root/.m2/repository/io/swagger/swagger-models/1.5.8/swagger-models-1.5.8.jar,/root/.m2/repository/io/swagger/swagger-annotations/1.5.8/swagger-annotations-1.5.8.jar,/root/.m2/repository/org/reflections/reflections/0.9.10/reflections-0.9.10.jar,/root/.m2/repository/com/google/code/findbugs/annotations/2.0.1/annotations-2.0.1.jar,/root/.m2/repository/com/fasterxml/jackson/jaxrs/jackson-jaxrs-json-provider/2.6.2/jackson-jaxrs-json-provider-2.6.2.jar,/root/.m2/repository/com/fasterxml/jackson/jaxrs/jackson-jaxrs-base/2.6.2/jackson-jaxrs-base-2.6.2.jar,/root/.m2/repository/org/glassfish/jersey/media/jersey-media-multipart/2.22.2/jersey-media-multipart-2.22.2.jar,/root/.m2/repository/org/jvnet/mimepull/mimepull/1.9.6/mimepull-1.9.6.jar,/root/data-lake-accelerator/commons/commons-spark/target/thinkbig-commons-spark-0.3.0-SNAPSHOT.jar,/root/data-lake-accelerator/core/schema-discovery/schema-discovery-model/target/thinkbig-schema-discovery-model-0.3.0-SNAPSHOT.jar,/root/data-lake-accelerator/services/feed-manager-service/feed-manager-api/target/thinkbig-feed-manager-api-0.3.0-SNAPSHOT.jar,/root/.m2/repository/org/springframework/spring-context-support/4.2.5.RELEASE/spring-context-support-4.2.5.RELEASE.jar,/root/.m2/repository/org/springframework/spring-core/4.2.5.RELEASE/spring-core-4.2.5.RELEASE.jar,/root/.m2/repository/com/google/guava/guava/19.0/guava-19.0.jar,/root/.m2/repository/org/apache/commons/commons-lang3/3.4/commons-lang3-3.4.jar,/root/.m2/repository/javax/inject/javax.inject/1/javax.inject-1.jar" --driver-java-options "-Xdebug -agentlib:jdwp=transport=dt_socket,server=y,suspend=n,address=5005 -Dcom.sun.management.jmxremote -Dcom.sun.management.jmxremote.port=9999 -Dcom.sun.management.jmxremote.authenticate=false -Dcom.sun.management.jmxremote.ssl=false -Dcom.sun.management.jmxremote.local.only=false -Djava.rmi.server.hostname=127.0.0.1" /root/data-lake-accelerator/services/spark-shell-service/target/thinkbig-spark-shell-service-0.3.0-SNAPSHOT.jar
